# deploy/supervisor.conf
[program:PROCESS_NAME]
directory=/var/www/apps/ENVT/FOLDER/
environment=REACT_APP_API_URL=https://API_URL_/v1/,NODEJS_PORT=PORT_
command=node server/app.js
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/DNS.error.log
stdout_logfile=/var/log/supervisor/DNS.out.log
